# app-emulation/xe-guest-utilities/files/85_xen-memory-hotplug.rules
SUBSYSTEM=="memory", ACTION=="add", RUN+="/bin/sh -c '[ -f /sys$devpath/state ] && echo online > /sys$devpath/state'"
